# Packer/elk/ansible/roles/curator/files/curator.service
[Unit]
Description=Run the curator script

[Service]
Type=simple
ExecStart=/home/elk/curator.sh

[Install]
WantedBy=multi-user.target
